# 2.Config/ansible-pull.crontab
*/30 * * * * ansible-pull -U https://github.com/CraigJPerry/home-network -d home-network -i 2.Config/hosts 2.Config/$(hostname).yml > /tmp/ansible-pull.$LOGNAME 2>&1

